# files/usr/lib/systemd/system/battery-full.target
[Unit]
Description=75% charged and attached to AC
DefaultDependencies=no
StopWhenUnneeded=yes
