# files/usr/lib/systemd/system/battery-full.target
[Unit]
Description=At least 70% of battery
DefaultDependencies=no
StopWhenUnneeded=yes
