# files/usr/lib/systemd/system/battery-low.target
[Unit]
Description=Battery Low (you can set the percentage in the udev rule)
DefaultDependencies=no
StopWhenUnneeded=yes
