# files/etc/systemd/system/rpm-ostreed-automatic.timer.d/override.conf
[Timer]
RandomizedDelaySec=10m
OnCalendar=*-*-* 4:00:00
Persistent=true